# docker/zato-web-admin/supervisord.conf
[supervisord]
nodaemon = true
logfile = /opt/zato/supervisord.log
loglevel = debug
user = root
directory = /opt/zato

[program:zato-web-admin]
command = /opt/zato/zato_start_web_admin
user = zato
