# docker/zato-server/supervisord.conf
[supervisord]
nodaemon = true
logfile = /opt/zato/supervisord.log
loglevel = debug
user = root
directory = /opt/zato

[program:zato-server]
command = /opt/zato/zato_start_server
user = zato
